# jollix/skripte/local
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Modified by Oliver Schwinn <patteh@web.de>
# and Jochen Spang <knochestolle@web.de>
# for jollix <www.jollix.de>
# location in cloop: /etc/init.d/local
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /home/xubuntu/berlios_backup/github/tmp-cvs/jollix/Repository/jollix/skripte/local,v 1.9 2003/04/27 16:25:08 knochen Exp $

HILITE="\033[1m"
NORMAL="\033[0m"

NET="yes"
CMDLINE="`cat /proc/cmdline`"
for x in $CMDLINE
do
	if [ "$x" = "nonet" ]
	then
		NET="no"
	fi
done

depend() {
	after *
}

blurb() {
	echo -ne ${HILITE}${1}
	#blurblen=${#1}
}

backup() {
	echo -ne "\033[0G\033[0K"
}

start() {
	ebegin "Freeing initrd"
	umount /initrd
	/sbin/freeramdisk /dev/ram0
	#make our mtab look sane...
	local mycdloc
	mycdloc=`cat /tmp/mycd`
	mycdloc=${mycdloc#/newroot*}
	mount -f $mycdloc /mnt/cdrom -t iso9660 -o ro
	if [ -e /mnt/livecd.cloop ]
	then
		mount -f /mnt/livecd.cloop /mnt/cloop -t ext2 -o ro
		#cloop was cached to tmpfs so that CD could be unmounted.
	else
		mount -f /mnt/cdrom/livecd.cloop /mnt/cloop -t ext2 -o ro
	fi
	rm /tmp/mycd
	mount -f tmpfs / -t tmpfs -o rw
	eend
	
        #modules-jollix, laden der Kernelmodule
        /sbin/modules-jollix
	#Partitionen mounten
	/sbin/transitmount
        /sbin/transitlink
}

stop() {
	return
}

# vim:ts=4
