# jollix/skripte/local
#!/sbin/runscript
# Copyright 1999-2002 Gentoo Technologies, Inc.
# Modified by Oliver Schwinn <patteh@web.de>
# and Jochen Spang <knochestolle@web.de>
# for jollix <www.jollix.de>
# location in cloop: /etc/init.d/local
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /home/xubuntu/berlios_backup/github/tmp-cvs/jollix/Repository/jollix/skripte/local,v 1.10 2003/09/08 09:50:06 patteh Exp $

HILITE="\033[1m"
NORMAL="\033[0m"
LANG_US="no"
ATIGAMER="no"
CMDLINE="`cat /proc/cmdline`"
for x in $CMDLINE
do
	if [ "$x" = "lang=us" ]
	then
		LANG_US="yes"
        elif [ "$x" = "atigamer" ]
        then
                ATIGAMER="yes"
	fi
done

depend() {
	after *
}

blurb() {
	echo -ne ${HILITE}${1}
	#blurblen=${#1}
}

backup() {
	echo -ne "\033[0G\033[0K"
}

start() {
	ebegin "Freeing initrd"
	umount /initrd
	/sbin/freeramdisk /dev/ram0
	#make our mtab look sane...
	local mycdloc
	mycdloc=`cat /tmp/mycd`
	mycdloc=${mycdloc#/newroot*}
	mount -f $mycdloc /mnt/cdrom -t iso9660 -o ro
	if [ -e /mnt/livecd.cloop ]
	then
		mount -f /mnt/livecd.cloop /mnt/cloop -t ext2 -o ro
		#cloop was cached to tmpfs so that CD could be unmounted.
	else
		mount -f /mnt/cdrom/livecd.cloop /mnt/cloop -t ext2 -o ro
	fi
	rm /tmp/mycd
	mount -f tmpfs / -t tmpfs -o rw
	eend
	
        #modules-jollix, laden der Kernelmodule
        /sbin/modules-jollix

	#Partitionen mounten
	/sbin/transitmount
        /sbin/transitlink

	#language
	if [ "$LANG_US" = "yes" ]
        then
            sed -i -e 's/de-latin1-nodeadkeys/us/' /etc/rc.conf
	    sed -i -e 's/"de"/"us"/' /etc/X11/XF86Config
	    /etc/init.d/keymaps start
            sed -i -e 's/Country=de/Country=us/' /home/user/.kde3.1/share/config/kdeglobals
	    sed -i -e 's/Language=de/Language=en_US/' /home/user/.kde3.1/share/config/kdeglobals
	    chown user /home/user/.kde3.1/share/config/kdeglobals
	    chgrp users /home/user/.kde3.1/share/config/kdeglobals
	else
	    /etc/init.d/keymaps start
        fi

        #ATI xinit workaround
        if [ "$ATIGAMER" = "yes" ]
        then
          touch /etc/atigamer
        fi 

        if [ ! -e /etc/atigamer ] ; then
          /etc/init.d/xdm start
        else
          touch /home/user/.startxinit0
          chown user /home/user/.startxinit0
          chgrp users /home/user/.startxinit0
        fi
}

stop() {
	return
}

# vim:ts=4
